# devops/app/supervisor/conf.d/cherry.conf
[program:cherry]
directory=/app
;pidfile=/app/tmp/cherry.pid
command=python3 cherry.py
user=app
autostart=true
autorestart=true
redirect_stderr=true
priority=800