# devops/app/supervisor/conf.d/migrate.conf
[program:migrate]
directory=/app
command=python3 manage.py migrate --noinput
user=app
redirect_stderr=true
priority=1